# gdm-background-switcher.service
[Unit]
Description=gdm-background-switcher
Before=gdm.service
[Service]
ExecStart=/root/gdm-switcher.sh
[Install]
WantedBy=multi-user.target
